# crontab
#
# crontab(5) for build
#
# TMP=/tmp EDITOR=emacs crontab crontab
#
# $Id$
#

SHELL=/bin/bash
MAILTO=build@lists.planet-lab.org

CVSROOT=:pserver:anon@cvs.planet-lab.org:/cvs
MODULE=build

#
# minute hour day-of-month month day-of-week
#

# install new crontab
0 0 * * *	cvs -d $CVSROOT checkout -p $MODULE/crontab 2>/dev/null | crontab -

# nightly build
17 5 * * *	SCRIPT=$(mktemp /tmp/build.sh.XXXXXX) ; cvs -d $CVSROOT checkout -p $MODULE/build.sh >$SCRIPT 2>/dev/null && $SHELL $SCRIPT -b $HOME/build/nightly -m $MAILTO -x 7 ; rm -f $SCRIPT
